# dserv-agent/dserv-agent.service
[Unit]
Description=dserv Management Agent
Documentation=https://github.com/SheinbergLab/dserv-agent
After=network.target
# Don't start before dserv, but don't require it
Wants=dserv.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/dserv-agent --registry https://dserv.net --workgroup YOUR-WORKGROUP-NAME --no-tls --allow-reboot --components /usr/local/dserv/config/components.json
Restart=always
RestartSec=5

# Security hardening (allow what we need)
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=read-only
PrivateTmp=true

# Environment file for token (optional)
EnvironmentFile=-/etc/dserv-agent/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dserv-agent

[Install]
WantedBy=multi-user.target
