# systemd/dserv.service
[Unit]
Description=Dataserver Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sheinb
ExecStart=/usr/local/dserv/dserv -c /usr/local/dserv/config/dsconf.tcl -t /usr/local/dserv/config/triggers.tcl
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target

